# node/supervisord.conf
[program:node-app]
command=npm start
directory=/app
autostart=true
autorestart=true

[program:cron]
command=service cron start
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
